# crontab
* * * * * python /mnt/c/dev/hello_docker_crontab/hello_docker.py
0-59 * * * * echo "hi" >> /root/cron.log
